# systemd/whattimenexttrain-frontend.service
[Unit]
Description=WhatTimeNextTrain Frontend Web Server
After=network.target whattimenexttrain-backend.service
Wants=network.target
Requires=whattimenexttrain-backend.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/WhatTimeNextTrain/frontend
Environment=NODE_ENV=production
ExecStartPre=/bin/sleep 15
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# セキュリティ設定
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/WhatTimeNextTrain/frontend/dist
PrivateTmp=true

[Install]
WantedBy=multi-user.target
